# etc/cron.d/dyndns-update
# /etc/cron.d/dyndns-update

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/2  *    * * *   root	/usr/local/sbin/dyndns-update check
