# etc/systemd/system/dyndns-update@.service
[Unit]
Description=DynDNS Update Script

[Service]
ExecStart=-/usr/local/sbin/dyndns-update xinetd
StandardInput=socket
User=root

